# i3/.config/autostart/polybar.desktop
[Desktop Entry]
Exec=/usr/bin/polybar -r i3
Name=polybar
Type=Application
OnlyShowIn=i3
