# i3/.config/autostart/polybar.desktop
[Desktop Entry]
Exec=/usr/bin/env bash -c "$HOME/.config/polybar/launch_polybar.sh"
Name=polybar
Type=Application
OnlyShowIn=i3
