# xinit/.config/autostart/synaptics_init.desktop
[Desktop Entry]
Exec=/home/josh/.local/lib/xinit/synaptics_init
Name=synaptics_init
Type=Application
NotShowIn=sway
