# xinit/.config/autostart/xcompmgr.desktop
[Desktop Entry]
Exec=/usr/bin/env xcompmgr
Name=xcompmgr
Comment=sample X compositing manager
Type=Application
NotShowIn=sway
