# xinit/.config/autostart/xrdb_init.desktop
[Desktop Entry]
Exec=/home/josh/.local/lib/xinit/xrdb_init
Name=xrdb_init
Type=Application
